# deploy/templates/systemd/vpn-worker.service.in
[Unit]
Description=VPN Control Panel Worker
After=network.target postgresql.service vpn-broker.service
Requires=vpn-broker.service

[Service]
Type=simple
User=${VCP_WORKER_USER}
Group=${VCP_WORKER_USER}
WorkingDirectory=${VCP_INSTALL_DIR}
Environment=PANEL_CONFIG_PATH=${VCP_CONFIG_DIR}/panel.yaml
Environment=VPN_SYSTEMCTL_CMD=sudo -n /bin/systemctl
ExecStart=${VCP_INSTALL_DIR}/.venv/bin/vpn-worker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
